# .github/package/prerm
#!/bin/sh

systemctl disable axonmq
systemctl stop axonmq
